# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/run/supervisord.pid

[program:php-fpm]
command=/usr/local/sbin/php-fpm -F
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:queue-worker]
; www-data must not inherit HOME=/root — libpq reads SSL certs from $HOME/.postgresql/
command=/bin/sh -c "export HOME=/var/www && exec php artisan queue:work --sleep=$${QUEUE_SLEEP:-3} --tries=$${QUEUE_TRIES:-3} --max-time=$${QUEUE_MAX_TIME:-3600} --timeout=$${QUEUE_TIMEOUT:-120}"
directory=/var/www/html
user=www-data
environment=HOME="/var/www",USER="www-data"
autostart=true
autorestart=true
priority=30
stopwaitsecs=130
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
